# netezza/power_failure_watchdog/pf_watchdog.cron
* * * * root /opt/pf_watchdog/pf_watchdog.sh >> /opt/pf_watchdog/cron.log